# rootfs/etc/services.d/snort/run
#!/usr/bin/execlineb -P

snort -Q -c /etc/snort/snort.lua --tweaks local -l /var/log/snort --create-pidfile
