# rootfs/etc/periodic/daily/10-pulledpork
#!/bin/sh

ulimit -t 60

set -e
/usr/bin/pulledpork.pl -P -c /etc/pulledpork/pulledpork.conf
cd /etc/snort/rules
/usr/bin/snort2lua -c /etc/snort/rules/snort2.rules -r /etc/snort/rules/snort.rules
